# build/php/files/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:phpfpm]
command=/run.sh
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true

[program:gearman]
command=/app/yii gearman/start %(process_num)s
process_name=worker-%(process_num)s
priority=1
numprocs=3
numprocs_start=1
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true

[program:export-gearman]
command=/app/yii gearman/start 10%(process_num)s
process_name=export-worker-10%(process_num)s
priority=2
numprocs=1
numprocs_start=1
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true

[program:hub-spot-gearman]
command=/app/yii gearman/start 20%(process_num)s
process_name=hub-spot-worker-20%(process_num)s
priority=3
numprocs=1
numprocs_start=1
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true

[program:dear-gearman]
command=/app/yii gearman/start 30%(process_num)s
process_name=dear-worker-30%(process_num)s
priority=3
numprocs=1
numprocs_start=1
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true

[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler